# ports-mgmt/caronade/files/caronade.in
#!/bin/sh
#
# PROVIDE: caronade
# REQUIRE: networking

. /etc/rc.subr

name="caronade"
rcvar="caronade_enable"

pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args="-P ${pidfile} -r -o /var/log/caronade.log %%PREFIX%%/bin/caronade -config %%PREFIX%%/etc/caronade/caronade.yaml"
command_chdir="%%PREFIX%%/caronade"

load_rc_config $name
: ${caronade_enable:=no}

run_rc_command "$1"
