# mail/smtp-proxy/files/smtpproxy.in
#!/bin/sh
#
# PROVIDE: smtpproxy
# REQUIRE: networking

. /etc/rc.subr

name="smtpproxy"
rcvar="smtpproxy_enable"

pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args="-P ${pidfile} -r -f %%LOCALBASE%%/bin/smtp-proxy -config %%ETCDIR%%/smtp-proxy.ini"

load_rc_config $name
: ${smtpproxy_enable:=no}

run_rc_command "$1"
